# .config/systemd/user/cgal-stop-docker-testsuite.timer
[Unit]
Description=Timer to stop CGAL Docker testsuite daily at 20:00

[Timer]
OnCalendar=20:00
Persistent=true

[Install]
WantedBy=timers.target
