# .config/systemd/user/podman.service.d/override.conf
[Service]
Environment=LOGGING="--log-level=info" "CONTAINER_CONF=/home/lrineau/.config/containers/no-remote.conf"
ExecStart=
ExecStart=/usr/bin/podman --remote=false $LOGGING system service
