# dot-config/systemd/user/cgal-testsuite.slice
[Unit]
Description=Slice for CGAL testsuite

[Slice]
# You can set resource limits here, for example:
#CPUQuota=50%
#MemoryMax=80%
