# dot-config/systemd/user/cgal-dump-crontab.timer
[Unit]
Description=Timer to dump crontab to Git every hour at minute 18

[Timer]
OnCalendar=*:18
Persistent=true

[Install]
WantedBy=timers.target
