# user/cgal-pull-testsuite-containers.timer

[Unit]
Description=Timer for pulling Docker images

[Timer]
OnCalendar=20:06

[Install]
WantedBy=timers.target