# user/cgal-stop-docker-testsuite.service
[Unit]
Description=Stop CGAL Docker testsuite

[Service]
Type=oneshot
ExecStart=/bin/echo "Stopping CGAL Docker testsuite"
